# samba_s6/etc/s6-overlay/s6-rc.d/samba/run
#!/command/execlineb -P
smbd --foreground --no-process-group --configfile=/config --debug-stdout